# consensus/metanode/deploy/metanode-go-sub@.service
[Unit]
Description=Metanode Go Sub (Node %i)
After=metanode-go-master@%i.service
Requires=metanode-go-master@%i.service

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{GO_WORKDIR}}/cmd/simple_chain
ExecStart={{GO_WORKDIR}}/cmd/simple_chain/simple_chain -config=config-sub-node%i.json

# Environment
Environment=GOTOOLCHAIN=go1.23.5
Environment=GOMEMLIMIT={{GOMEMLIMIT}}
Environment=XAPIAN_BASE_PATH=sample/node%i/data/data_sub/xapian_node

# File descriptors
LimitNOFILE=100000

# Auto-restart
Restart=always
RestartSec=3
StartLimitIntervalSec=300
StartLimitBurst=10

# Logging
StandardOutput=append:{{RUST_WORKDIR}}/logs/node_%i/go-sub-stdout.log
StandardError=append:{{RUST_WORKDIR}}/logs/node_%i/go-sub-stdout.log

[Install]
WantedBy=multi-user.target
